# petehome.service
[Unit]
Description=Petehome Smart Home Control Center
After=network.target

[Service]
Type=simple
User=pete
WorkingDirectory=/home/pete/dev/petehome
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
ExecStart=/usr/bin/pm2 start ecosystem.config.js --no-daemon
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
